# dist/bitwarden-keyring.service
[Unit]
Description=Bitwarden Secret Service
Documentation=https://github.com/joe/bitwarden-keyring
Before=default.target

[Service]
Type=dbus
BusName=org.freedesktop.secrets
ExecStart=bitwarden-keyring
Restart=on-failure
RestartSec=5
RuntimeDirectory=bitwarden-keyring

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths="%h/.config/Bitwarden CLI"
PrivateTmp=yes

[Install]
WantedBy=default.target
